# gpu-board-server.service
[Unit]
Description=Gpu Board Server
After=network.target

[Service]
ExecStart=
